# common/apps/mgr/sbin/edge-mgr.timer
[Unit]
Description=timer for edge-mgr service

[Timer]
OnBootSec=1min

[Install]
WantedBy=timers.target
